# workstation/linux/systemd/system/aidecheck.timer
[Unit]
Description=Daily AIDE Integrity Check

[Timer]
OnCalendar=daily

Persistent=true

Unit=aidecheck.service

[Install]
WantedBy=timers.target
