# workstation/linux/systemd/system/lynis.timer
[Unit]
Description=Weekly Lynis Security Audit

[Timer]
OnCalendar=weekly

Persistent=true

Unit=lynis.service

[Install]
WantedBy=timers.target
